# conf/supervisord.conf
[group:drawnby]
programs=drawnby_app

[program:drawnby_app]
command=python manage.py runserver
directory=/home/dash/drawnby/
user=nobody
autostart=true
autorestart=true

