# conf/supervisord.conf
[program:drawnby_app]
command=python manage.py runserver 0.0.0.0:80 --noreload
directory=/home/dash/drawnby/
user=nobody
autostart=true
autorestart=true
